# systemd/shibe-upgrade.service
[Unit]
Description=Automatic Shibe OS Upgrade Service

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc upgrade
