# systemd/shibe-upgrade.timer
[Unit]
Description=Run Shibe OS Upgrade Daily

[Timer]
OnCalendar=03:00
Persistent=true

[Install]
WantedBy=timers.target
